# packaging/dockguard.service
[Unit]
Description=dockguard - Lightweight Docker container update watcher
Documentation=https://github.com/h4llow3En/dockguard
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/bin/dockguard
Restart=on-failure
RestartSec=5

# Set options via environment variables if desired:
# Environment="GUARD_LOG_LEVEL=info"
# Environment="GUARD_CLEAN=false"

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
